# debian/prerm
#!/bin/bash
systemctl stop wildfly.service || true
systemctl disable wildfly.service || true
